# backend/hotel-backend.service
[Unit]
Description=Hotel Management Backend API
Documentation=https://github.com/nvdong123/zalo_be
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=exec
User=hotelapp
Group=hotelapp
WorkingDirectory=/var/www/hotel-backend/backend
Environment=PATH=/var/www/hotel-backend/backend/venv/bin
Environment=PYTHONPATH=/var/www/hotel-backend/backend
ExecStart=/var/www/hotel-backend/backend/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --timeout 120 \
    --keepalive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --access-logfile /var/log/hotel-backend/access.log \
    --error-logfile /var/log/hotel-backend/error.log \
    --log-level info \
    app.main:app

Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hotel-backend

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/uploads /var/log/hotel-backend /var/www/hotel-backend/backend/logs
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
